# utils/etc/profile.d/hf.sh
export HF_HUB_DISABLE_TELEMETRY=1
export HF_HUB_ENABLE_HF_TRANSFER=1
export HF_DATASETS_CACHE=/data/LLM/huggingface/datasets
export HF_HUB_CACHE=/data/LLM/huggingface/hub
export HF_XET_CACHE=/data/LLM/huggingface/xet
export HF_ASSETS_CACHE=/data/LLM/huggingface/assets
export NLTK_DATA=/data/LLM/nltk_data
