# utils/etc/profile.d/nvidia-mps.sh
export CUDA_MPS_PIPE_DIRECTORY="/tmp/nvidia-mps-$USER"
export CUDA_MPS_LOG_DIRECTORY="/tmp/nvidia-log-$USER"
